# php-worker/supervisord.d/admin.ch.conf
[program:admin-ch]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ch/admin/artisan horizon
autostart=true
autorestart=true
numprocs=8
user=laradock
redirect_stderr=true
